# mariadb-repl.service
[Unit]
Description=MySQL Stack Deployment
After=network.target

[Service]
ExecStart=/etc/init.d/mariadb-start.sh

[Install]
WantedBy=default.target