# mariadb-repl.service
[Unit]
Description=MySQL Stack Deployment
After=network.target

[Service]
ExecStart=/var/lib/mariadb/start.sh

[Install]
WantedBy=default.target